# templates/systemd/vpn-proxy.socket
[Socket]
ListenStream=1081
[Install]
WantedBy=sockets.target
