# templates/systemd/vpn-proxy.socket
[Socket]
ListenStream=1081
BindToDevice=lo
[Install]
WantedBy=sockets.target
